# ops/com.victor.ai-job-search.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.victor.ai-job-search.digest</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>/Users/victorblanco/Documents/ai-job-search/tools/daily_digest.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/victorblanco/Documents/ai-job-search</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>RunAtLoad</key>
	<false/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>StandardOutPath</key>
	<string>/Users/victorblanco/Library/Logs/ai-job-search-digest.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/victorblanco/Library/Logs/ai-job-search-digest.error.log</string>
</dict>
</plist>
